# deploy/systemd/ucc-bulk.service
[Unit]
Description=UCC Ingestion — State Bulk Subscriptions (15 states)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ucc
Group=ucc
WorkingDirectory=/opt/ucc/ucc-registration
ExecStart=/usr/bin/python3 /opt/ucc/deploy/fetch_secrets.py run --tier state_bulk
TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ucc-bulk
